# meta-axiado/meta-ax3000/recipes-ax3000/ax-apps/ax-utils/logmgr/ax_logger.service
[Unit]
Description=TCU Logger Application Service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
ExecStart=/usr/bin/ax_logger
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
