# meta-axiado/meta-ax3000/recipes-ax3000/ax-apps/ax-utils/flowmeter/flowmeter.service
[Unit]
Description=Flowmeter Service
After=eip_fwl.service

[Service]
ExecStart=/usr/bin/flowmeter
Restart=always

[Install]
WantedBy=multi-user.target