# meta-axiado/meta-ax3000/recipes-ax3000/ax-apps/ax-utils/flowmeter/rwui.service
[Unit]
Description=Ransomware Service
After=flowmeter.service

[Service]
ExecStart=/usr/bin/ransomware-predictor
Restart=always

[Install]
WantedBy=multi-user.target