# meta-axiado/meta-ax3000/recipes-ax3000/ax-apps/ax-agent/demo.service
[Unit]
Description=demo bmc agent
After=dbus.service

[Service]
ExecStart=/usr/bin/axdemo
Restart=always

[Install]
WantedBy=multi-user.target
