# meta-axiado/meta-ax3000/recipes-ax3000/ax-apps/ax-agent/test-data.service
[Unit]
Description=Random test data for bmc agent
After=dbus.service demo.service

[Service]
ExecStart=/usr/bin/ax-demo-test
Restart=always

[Install]
WantedBy=multi-user.target
