# meta-axiado/meta-ax3000/recipes-ax3000/ax-apps/ax-utils/sysproxy/sysproxy.service
[Unit]
Description=Sysproxy Service

[Service]
ExecStart=/usr/bin/sysmgr_proxy
Restart=always

[Install]
WantedBy=multi-user.target
